# MSS5004W/OpenWrt-MSS5004W-Dev/vendor/uci-defaults/vendor-settings
#!/bin/sh

hostname=""
portmap=""

if [ -f  /etc/config/voip ] ; then
	hostname="MSS5003WV"
	portmap="lan1 lan2 lan3 wlan1 wlan2"
else
	hostname="MSS5004W"
	portmap="lan1 lan2 lan3 lan4 wlan1 wlan2"
	uci batch <<-EOF
		set wireless.ra1.enabled=0
		commit wireless
	EOF
fi

uci -q batch <<-EOF >/dev/null
	delete wanctl.wanlink_0
	commit wanctl
	add wanctl wanlink
	set wanctl.wanlink_0.VlanID=-1
	set wanctl.wanlink_0.Cos=0
	set wanctl.wanlink_0.Servicemode=0
	set wanctl.wanlink_0.Mode=0
	set wanctl.wanlink_0.PortMap="$portmap"
	commit wanctl
	set system.@system[0].hostname="$hostname"
	commit system
EOF
	
echo "/etc/config/wanctl" >/usr/lib/opkg/keep.conf
echo "$hostname" > /proc/sys/kernel/hostname
sed -i 's/\(DeviceName=\)[^\n]*/\1'"${hostname}"'/g' /tmp/MME_devinfo
exit 0
