# MSS5004W/OpenWrt-MSS5004W-Dev/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
ttyS1::respawn:/sbin/getty -L ttyS1 57600 vt100
